# deployment/etc/cron.hourly/sniffer-append-to-log
#!/bin/sh
echo  `date` >> /var/log/aircrack-ng/cronlog.log
python /home/pi/sniffer/append2log.py >> /var/log/aircrack-ng/cronlog.log

