# deployment/etc/cron.hourly/sniffer-append-to-log
#!/bin/sh
echo  `date` >> /home/pi/sniffer/cronlog.log

python /home/pi/sniffer/append2log.py >> /home/pi/sniffer/cronlog.log

